# infra/templates/learning-agent.service.j2
[Unit]
Description=Learning Agent
[Service]
ExecStart=/usr/bin/java -jar /home/ubuntu/Learning.jar
Type=simple
User=ubuntu
[Install]
WantedBy=multi-user.target